# systemd/algo-bot-monitor.service
[Unit]
Description=Algo-bot position monitor — polls Alpaca fills and closes journal entries
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/madanw/projects/algo-bot
EnvironmentFile=/home/madanw/projects/algo-bot/.env
ExecStart=/home/madanw/projects/algo-bot/.venv/bin/python -m backend.scripts.monitor
Restart=on-failure
RestartSec=60
StandardOutput=append:/home/madanw/algo-bot-monitor.log
StandardError=append:/home/madanw/algo-bot-monitor.log

[Install]
WantedBy=default.target
