# systemd/algo-bot.timer
[Unit]
Description=Fire algo-bot weekdays 09:45 America/New_York

[Timer]
OnCalendar=Mon..Fri 09:45 America/New_York
Persistent=false
Unit=algo-bot.service

[Install]
WantedBy=timers.target
